# etc/cron.d/nagios-plugins-cyb
PATH=/bin:/usr/bin:/sbin:/usr/sbin

10 * * * * root if [ -x /usr/sbin/gnt-cluster ] && ps -C ganeti-masterd 2>&1 >/dev/null; then /usr/sbin/cron_check_ganeti_cluster_verify; fi
*/5 * * * * nagios if [ -x /usr/bin/cron_check_rbl ]; then /usr/bin/cron_check_rbl; fi
