# debian/cron.d
*/15 * * * * root nice -20 ionice -c3 /usr/bin/check_safe_browsing_bulk
