# debian/cron.d
*/15 * * * * root nice -20 ionice -c3 /usr/bin/nagios-plugins-cyb-update-gsb-db
